# kubeproxy.service
[Unit]
Description=Proxy for kube
[Service]
ExecStart=/home/ubuntu/kube_proxy_start.sh
Restart=always
[Install]
WantedBy=multi-user.target
